# com.mac.scraper.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.mac.scraper.hourly</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/mac/Documents/scraper/venv/bin/python</string>
		<string>/Users/mac/Documents/scraper/run_all_scrapers.py</string>
	</array>

	<key>StartInterval</key>
	<integer>3600</integer>

	<key>WorkingDirectory</key>
	<string>/Users/mac/Documents/scraper</string>

	<key>StandardOutPath</key>
	<string>/Users/mac/Documents/scraper/scraper_hourly.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/mac/Documents/scraper/scraper_hourly.err.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
	</dict>
</dict>
</plist>
